# thumbai.service
# GENERATED BY aah CLI - Feel free to customization it.
# FILE: thumbai.service
# DATE: Sun, 18 Nov 2018 03:37:34 -0800
# DESC: aah application systemd service file

[Unit]
Description=thumbai application
After=network.target

[Service]
#User=thumbai
#Group=thumbai
#EnvironmentFile=/home/app/thumbai_env_values
ExecStart=/home/app/thumbai run --config /home/app/thumbai.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
